# root/etc/services.d/gphotos-sync/run
#!/usr/bin/with-contenv bash

# exec gphotos-sync /storage
